# systemd/sysgit_commit@.path
[Path]
PathChanged=%I
TriggerLimitIntervalSec=30min

[Install]
WantedBy=paths.target
